# linux-f10-dbus/pkg-plist
bin/dbus-cleanup-sockets
bin/dbus-daemon
bin/dbus-monitor
bin/dbus-send
bin/dbus-uuidgen
etc/dbus-1/session.conf
etc/dbus-1/system.conf
etc/rc.d/init.d/messagebus
lib/dbus-1/dbus-daemon-launch-helper
usr/share/doc/dbus-1.2.4/COPYING
usr/share/man/man1/dbus-cleanup-sockets.1.gz
usr/share/man/man1/dbus-daemon.1.gz
usr/share/man/man1/dbus-monitor.1.gz
usr/share/man/man1/dbus-send.1.gz
usr/share/man/man1/dbus-uuidgen.1.gz
@dirrm var/run/dbus
@dirrm var/lib/dbus
@dirrm usr/share/doc/dbus-1.2.4
@dirrm usr/share/dbus-1/system-services
@dirrm usr/share/dbus-1/services
@dirrm usr/share/dbus-1/interfaces
@dirrm usr/share/dbus-1
@dirrm lib/dbus-1
@dirrm etc/dbus-1/system.d
@dirrm etc/dbus-1/session.d
@dirrm etc/dbus-1
